# 90-backlight.rules
SUBSYSTEM=="backlight", ACTION=="add", RUN+="/usr/bin/chgrp video %S%p/brightness", RUN+="/usr/bin/chmod g+w %S%p/brightness"
